# mysql.supervisor.conf
[program:mysql]
command=mysqld --datadir=/var/lib/mysql --user=mysql
user=mysql
